# util/socket-activation/endlessh.socket
[Socket]
ListenStream=[::]:22

[Install]
WantedBy=sockets.target
